# scripts/base/etc/rc.local
#!/bin/sh -e

# FIXME https://www.troublenow.org/752/debian-10-add-rc-local/
# Power-recovery

echo 37 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio37/direction
echo 1 > /sys/class/gpio/gpio37/value
