# scripts/base/etc/cron.monthly/autofsck
#!/bin/bash -ex

date > /var/log/autofsck.log
/usr/sbin/nas-fsck &>> /var/log/autofsck.log
